# build/tkfs.service
[Unit]
Description=TrustedKeep Encrypted Filesystem
After=syslog.target network.target

[Service]
LimitNOFILE=32768
LimitMEMLOCK=infinity
ExecStart=/usr/local/tkfs/bin/tkfs -c /usr/local/tkfs/config/sample_config.yaml
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
